# cron/accounts.cron
0 * * * *  cd /home/jacob/public_html/yase/cron/ && php daemon.php jacob >> jacob.log
